# scripts/deposit_import.crontab
*/5 * * * * /libra2/scripts/deposit_import.sh
